# systemd/.config/systemd/user/mp_rater.service
[Unit]
Description=watches for playcounts and skips
After=mpd.service mako.service

[Service]
Type=exec
ExecStart=%h/.local/bin/mp_rater --use-tags -vvvv -L
Restart=always

[Install]
WantedBy=default.target
