# systemd/.config/systemd/user/test.service
[Unit]
Description=syncs the mails 

[Service]
Type=simple
ExecStart=/home/gireesh/Trash/test.sh

[Install]
WantedBy=default.target
