# systemd/.config/systemd/user/wallpaper_set.service
[Unit]
Description=sets the background based on the time
