# configs/systemd/user/mp_rater.service.d/override.conf
[Unit]
Description=watches for playcounts and skips
After=mpd.service

[Service]

[Install]
WantedBy=graphical.target