# homemanager-modules/configs/systemd/user/mail_sync.service
[Unit]
Description=syncs the mails 

[Service]
Type=exec
ExecStart=/home/gireesh/.local/bin/mailsync

[Install]
WantedBy=default.target
