# homemanager-modules/configs/systemd/user/wallpaper_set.service.d/override.conf
[Unit]
Description=sets the background based on the time

[Service]
ExecStart=/home/gireesh/.config/sway/scripts/bgset.sh

[Install]
WantedBy=graphical.target